# crontab
0 18 * * * /Users/jhatfield/scripts/update_ssh_config
@reboot    /Users/jhatfield/scripts/update_ssh_config
